# crontab
* * * * * /usr/bin/php /My/home/hue.php >/dev/null 2>&1
